# proton-sync.path
[Unit]
Description=Proton Drive Path Watcher

[Path]
PathChanged=%h/Documents
PathChanged=%h/Pictures
PathChanged=%h/Videos
TriggerLimitIntervalSec=60

[Install]
WantedBy=default.target
